# deploy/vps/systemd/discord-bot-poll.timer
[Unit]
Description=Poll GitHub main every 5 minutes (webhook fallback)

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
